# ai-newsletter/macos/com.kakao.ai-newsletter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kakao.ai-newsletter</string>

  <key>WorkingDirectory</key>
  <string>/Users/kakao/aiprojects/aiexpert/ai-newsletter</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/kakao/aiprojects/aiexpert/ai-newsletter/.venv/bin/python</string>
    <string>main.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>11</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/kakao/aiprojects/aiexpert/ai-newsletter/newsletter.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kakao/aiprojects/aiexpert/ai-newsletter/newsletter.err.log</string>
</dict>
</plist>
